# templates/vm-systemd.socket.j2
[Unit]
Description=QEMU KVM {{ kvm_machine.name }} spice socket

[Socket]
ListenStream=[{{ kvm_machine_spice_addr }}]:{{ kvm_machine_spice_port }}

[Install]
WantedBy=sockets.target
